# src/debian/ddns-updater.cron.d
## Start DDNS Updater every 10 minutes
*/10 * * * * root /usr/sbin/ddns-updater > /dev/null 2>&1
