# src/debian/ddns-updater/etc/cron.d/ddns-updater
## Start ipcheck DDNS updater every 10 minutes
*/10 * * * * root /usr/sbin/ddns-updater
